# platform/ciscovs/docker-syncd-ciscovs/supervisord.conf
[supervisord]
logfile_maxbytes=1MB
logfile_backups=2
nodaemon=true

[eventlistener:supervisor-proc-exit-listener]
command=/usr/bin/supervisor-proc-exit-listener --container-name syncd
events=PROCESS_STATE_EXITED
autostart=true
autorestart=unexpected

[program:start.sh]
command=/usr/bin/start.sh
priority=1
autostart=true
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
priority=2
autostart=false
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog

[program:veth-creator]
command=/etc/scripts/veth-create.sh
priority=2
startsecs=0
autostart=true
autorestart=false
dependent_startup=true

[program:tc-creator]
command=/etc/scripts/tc-create.sh
priority=2
startsecs=0
autostart=true
autorestart=false
dependent_startup=true
dependent_startup_wait_for=veth-creator:exited

[program:set-promiscuous]
command=/etc/scripts/set-promisc.sh
priority=2
startsecs=0
autostart=true
autorestart=false
dependent_startup=true
dependent_startup_wait_for=tc-creator:exited

[program:nsim]
command=/bin/bash -c 'cd /etc/scripts/ && python3 ./nsim_inst.py'
priority=5
autostart=true
autorestart=false
dependent_startup=true
dependent_startup_wait_for=rsyslogd:running veth-creator:exited

[program:nsim_kernel]
command=/bin/bash -c 'python3 /etc/scripts/nsim_init.py'
priority=5
autostart=true
autorestart=true
dependent_startup=true

[program:syncd]
command=/usr/bin/syncd_start.sh
priority=5
autostart=false
autorestart=false
dependent_startup=true
dependent_startup_wait_for=rsyslogd:running veth-creator:exited nsim:running nsim_kernel:running


[eventlistener:syncd-watchdog]
command=/usr/bin/app-watchdog -d syncd
events=PROCESS_STATE,TICK_60
stdout_logfile=stdout
stderr_logfile=stderr

